# deploy/systemd/nexusdesk-openclaw-events.service
[Unit]
Description=NexusDesk OpenClaw Event Daemon
After=network.target

[Service]
WorkingDirectory=/opt/nexusdesk/backend
EnvironmentFile=/opt/nexusdesk/backend/.env.production
ExecStart=/opt/nexusdesk/venv/bin/python scripts/run_openclaw_event_daemon.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
