# deploy/systemd/nexusdesk-worker.service
[Unit]
Description=NexusDesk Worker
After=network.target nexusdesk-openclaw-bridge.service
Wants=nexusdesk-openclaw-bridge.service

[Service]
WorkingDirectory=/opt/nexusdesk/backend
EnvironmentFile=/opt/nexusdesk/backend/.env.production
ExecStartPre=/opt/nexusdesk/venv/bin/python /opt/nexusdesk/backend/scripts/wait_for_openclaw_bridge.py --timeout 20
ExecStart=/opt/nexusdesk/venv/bin/python scripts/run_worker.py --worker-id worker-main
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
